# src/auto-hibernate.cron
# try to hibernate every 5 minutes
*/5 *  * * *	root	/usr/sbin/auto-hibernate.sh

